# scripts/com.socialmediaempire.videopipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.socialmediaempire.videopipeline</string>

    <!--
      Use system python3 (/usr/bin/python3) — the repo's requirements are
      installed against it, and video_automation/__init__.py eagerly pulls
      in google.genai and other heavy deps that aren't on the homebrew
      pythons. If you switch pythons, reinstall requirements.txt first.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/homefolder/Desktop/social-media-empire/scripts/local_video_pipeline.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/homefolder/Desktop/social-media-empire</string>

    <!--
      Defensive env: launchd does not inherit a login shell's locale, so
      Python 3.9 can fail with `init_fs_encoding` errors on some macOS
      configurations. Pin UTF-8 explicitly. (Note: this does NOT bypass
      Full Disk Access — see docs/VIDEO_PIPELINE_TROUBLESHOOT.md if you
      see `[Errno 1] Operation not permitted` in the error log.)
    -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONIOENCODING</key>
        <string>utf-8</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>LC_ALL</key>
        <string>en_US.UTF-8</string>
    </dict>

    <!-- Every 8 hours, matching the content-engine cron cadence. -->
    <key>StartInterval</key>
    <integer>28800</integer>

    <!-- Run once at login so the first render happens without waiting 8h. -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Be quiet about CPU if the system is busy — renders aren't urgent. -->
    <key>Nice</key>
    <integer>5</integer>

    <key>StandardOutPath</key>
    <string>/Users/homefolder/Desktop/social-media-empire/logs/video_pipeline.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/homefolder/Desktop/social-media-empire/logs/video_pipeline_error.log</string>
</dict>
</plist>
